# harvester/harvester/supervisord/celery.conf
[program:celery]
directory=/usr/src/app
command=celery -A harvester worker
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
